# cron.d/webmon
*/5 * * * * root /usr/local/bin/webmon.sh >/dev/null 2>&1
